# pi/rc.local
# Earthstar
/home/pi/earthstar-control-software/pi/effectbox.sh &
/home/pi/earthstar-control-software/pi/led_driver.sh &
